# ops/launchd/com.finsider.symphony.mitch-be.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.finsider.symphony.mitch-be</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>exec "$HOME/Projects/symphony/ops/start-symphony.sh" "$HOME/Projects/symphony/workflows/WORKFLOW.mitch-be.md" 4311</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>30</integer>

  <key>StandardOutPath</key>
  <string>/Users/dm3n/Projects/symphony/logs/launchd.mitch-be.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/dm3n/Projects/symphony/logs/launchd.mitch-be.err.log</string>
</dict>
</plist>
